# conf/supervisor.ini
[program:ndeploy_watcher]
command=/opt/nDeploy/scripts/watcher.py start
autorestart=true
