# rpm_buildtree/nDeploy-pkg/etc/cron.d/ndeploy
MAILTO=""
*/5 * * * * root /opt/nDeploy/scripts/monitor_apache.php >/dev/null
*/5 * * * * root /opt/nDeploy/scripts/monitor_nginx_watcher.php >/dev/null
