# rpm_buildtree/nDeploy-pkg/etc/cron.d/ndeploy
MAILTO=""
*/2 * * * * root /root/monitor_apache.php >/dev/null
*/2 * * * * root /opt/nDeploy/scripts/monitor_php_fpm.pl >/dev/null