# rpm_buildtree/nDeploy-pkg/etc/cron.d/ndeploy-phpbrew
# 0 * * * * /etc/init.d/ndeploy_backends restart
* * * * * /opt/nDeploy/scripts/monitor_php_fpm.sh >/dev/null 2>&1